# imageroot/systemd/user/kickstart.service
#
# Copyright (C) 2022 Nethesis S.r.l.
# SPDX-License-Identifier: GPL-3.0-or-later
#

#
# This systemd unit starts a zusam instance using Podman.
# Most parts of this file come from podman-generate-systemd.
#

[Unit]
Description=Podman zusam.service
Requires=mariadb-app.service zusam-app.service
Before=mariadb-app.service zusam-app.service

[Service]
Environment=PODMAN_SYSTEMD_UNIT=%n
EnvironmentFile=-%S/state/environment
Restart=always
TimeoutStopSec=70
ExecStartPre=/bin/rm -f %t/zusam.pid %t/zusam.pod-id
ExecStartPre=/usr/bin/podman pod create --infra-conmon-pidfile %t/zusam.pid \
    --pod-id-file %t/zusam.pod-id \
    --name zusam \
    --publish 127.0.0.1:${TCP_PORT}:80 \
    --replace
ExecStart=/usr/bin/podman pod start --pod-id-file %t/zusam.pod-id
ExecStop=/usr/bin/podman pod stop --ignore --pod-id-file %t/zusam.pod-id -t 10
ExecStopPost=/usr/bin/podman pod rm --ignore -f --pod-id-file %t/zusam.pod-id
PIDFile=%t/zusam.pid
Type=forking

[Install]
WantedBy=default.target
